# templates/init/systemd/with-foreground-irqbalance.service.erb
# This file was generated by Puppet

[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=<%= @config_dir_path -%>/<%= @config_file_name -%>
ExecStart=<%= @irqbalance_path -%> --foreground $IRQBALANCE_ARGS

[Install]
WantedBy=multi-user.target
