# templates/init/systemd/with-foreground-irqbalance.service.erb
# This file was generated by Puppet
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=<%= @config %>
ExecStart=/usr/sbin/irqbalance --foreground $IRQBALANCE_ARGS

[Install]
WantedBy=multi-user.target
