# templates/init/systemd/without-foreground-irqbalance.service.erb
# This file was generated by Puppet
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=<%= @config %>
Type=forking
ExecStart=/usr/sbin/irqbalance $ONESHOT

[Install]
WantedBy=multi-user.target
