# templates/init/systemd/without-foreground-irqbalance.service.erb
# This file was generated by Puppet
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=<%= @config_dir_path -%>/<%= @config_file_name -%>
Type=forking
ExecStart=<%= @irqbalance_path -%> $ONESHOT

[Install]
WantedBy=multi-user.target
