# templates/init/upstart/debian-noargs-irqbalance.conf.erb
# This file was generated by Puppet

description	"CPU interrupts balancing daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork

script
	test -f <%= @config_dir_path -%>/<%= @config_file_name -%> && . <%= @config_dir_path -%>/<%= @config_file_name -%>

	test "$ENABLED" != "0" || exit 0

	if test "$ONESHOT" != "0"; then
		DOPTIONS="--oneshot"
	fi

	exec <%= @irqbalance_path -%>  $DOPTIONS

end script
