# templates/init/upstart/el6-irqbalance.conf.erb
# This file was generated by Puppet
# It is a modification of the Upstart script provided by Ubuntu
# start and stop irqbalance daemon

# The irqbalance daemon will distribute interrupts across
# the cpus on a multiprocessor system with the purpose of
# spreading the load

description	"IRQ balancing daemon"

start on runlevel [2345]
stop on runlevel [!2345]

script 
	test -f <%= @config %> && . <%= @config %>

	test "$ENABLED" != "0" || exit 0

	if test "$ONESHOT" != "0"; then
		DOPTIONS="--oneshot"
	fi

	exec /usr/sbin/irqbalance  $DOPTIONS $IRQBALANCE_ARGS --foreground

end script
